# inst/beeco.crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin


0 */5 * * * Rscript --default-packages=apis -e ' kineis_save_real_time();'
